# wsl/etc/profile.d/proxy.sh
export HTTP_PROXY=http://172.29.128.1:7897
export HTTPS_PROXY=http://172.29.128.1:7897
export NO_PROXY="localhost,127.0.0.1"
